# projects/saas/queue/tasks.cron
*/5 * * * * * http://localhost:8000/url/to/do/something